# musicbot.hardened.service
[Unit]
Description=Discord Music Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/discord
EnvironmentFile=/etc/discord-music.env
Environment=PATH=/home/ubuntu/discord/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStartPre=/usr/bin/bash -lc 'command -v ffmpeg >/dev/null || { echo "ffmpeg not found"; exit 1; }'
ExecStart=/home/ubuntu/discord/.venv/bin/python /home/ubuntu/discord/botmusic.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
